# templates/default/upstart.conf.erb
description "statsd"
author      "dan@shopkeep.com"

start on startup
stop on shutdown

script
    # We found $HOME is needed. Without it, we ran into problems
    export HOME="/root"
    exec sudo -u <%= @user %> node <%= @basedir %>/stats.js <%= @confdir %>/localConfig.js 2>&1 >> <%= @logdir %>/statsd.log
end script

post-start script
	PID=`status statsd | egrep -oi '([0-9]+)$' | head -n1`
	echo $PID > <%= @pid_file %>
end script

post-stop script
	rm -f <%= @pid_file %>
end script
